# contrib/sigaw-daemon.service.in
[Unit]
Description=Sigaw - Discord voice overlay daemon
Documentation=https://git.macco.dev/macco/sigaw
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=@sigaw_daemon_path@ --foreground
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/dev/shm
ReadWritePaths=%t
ReadWritePaths=%h/.config/sigaw
ReadWritePaths=%h/.cache/sigaw

[Install]
WantedBy=graphical-session.target
